# deploy/water-quality.service
[Unit]
Description=美兰机场供水站水质管理系统
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/water-quality/backend
EnvironmentFile=/opt/water-quality/.env
ExecStart=/opt/water-quality/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000
Restart=on-failure
RestartSec=5

# 安全加固
NoNewPrivileges=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
